# .config/autostart/keymaps.desktop
[Desktop Entry]
Name=Keymaps
Comment=Sets the user keymaps
Exec=bash keymapping.sh
Terminal=false
Type=Application
